# archiso/airootfs/etc/xdg/autostart/ksuperkey.desktop
[Desktop Entry]
Name=ksuperkey
Exec=ksuperkey -e 'Super_L=Alt_L|F1'
Terminal=false
StartupNotify=false
